# ruby/ge-svc/ge-svc.service
[Unit]
Description=Grand Exchange API Proxy Service

[Service]
ExecStart=/bin/bash -l -c ge-svc.rb
User=rsmm

[Install]
WantedBy=multi-user.target
