# ruby/ge-svc/ge-svc.service
[Unit]
Description=Grand Exchange API Proxy Service

[Service]
ExecStart=/usr/bin/ge-svc.rb
User=rsmm

[Install]
WantedBy=multi-user.target
